# privacy-anonymizer/systemd/privacy-anonymizer.service
[Unit]
Description=Anonimizador de Documentos (Privacy Anonymizer)
After=network.target

[Service]
Type=simple
User=server-paulo
WorkingDirectory=/home/server-paulo/Software/privacy-openai/privacy-anonymizer
Environment="PATH=/home/server-paulo/.local/bin:/usr/bin:%PATH%"
Environment="OPF_DEVICE=cpu"
Environment="FLASK_HOST=0.0.0.0"
Environment="FLASK_PORT=5042"
ExecStart=/home/server-paulo/.local/bin/gunicorn --bind 0.0.0.0:5042 --workers 1 --threads 4 --timeout 300 --access-logfile - --error-logfile - wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
